# files/argononed-rc-run
#!/sbin/openrc-run

name="argononed"
command="/usr/local/bin/argononed"
command_args=""
pidfile="/run/argononed.pid"
command_background=true
respawn_delay=5
respawn_max=10

stop() {
	/usr/local/bin/argononed-poweroff
}
